# server/messenger.service
[Unit]
Description=Messenger API Server
After=network.target

[Service]
User=root
# Change this to your actual project directory
WorkingDirectory=/root
# Ensure python3 and uvicorn paths are correct
ExecStart=/usr/bin/python3 -m uvicorn server:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
